# crontab
10,50 12 * * * cd ~/psg; python kuangzanThread.py -e
30 12 * * * cd ~/psg; python soukuangThread.py

11 19 * * * cd ~/psg; python soukuangThread.py
10 21 * * * cd ~/psg; python soukuangThread.py
2 4,12,23 * * * cd ~/psg; python oneTimePerDay.py

16 21 * * * cd ~/psg; python tongsang.py
17 * * * * cd ~/psg; sh restart_tax.sh
28 */2 * * * cd ~/psg; sh restart_touzi.sh
23 */3 * * * cd ~/psg; python qiangzengThread.py -t 1
